# docker/idle.service
[Unit]
Description=Idle controller service

[Service]
ExecStart=/usr/local/bin/stop-idle-services
PrivateTmp=true
RemainAfterExit=false
