# etc/systemd/system/fatc.service
[Unit]
Description=Controller service
Requires=configuration.service idle.timer
After=configuration.service idle.timer

[Service]
ExecStart=/usr/bin/fatctl start
ExecStop=/usr/bin/fatctl stop
PrivateTmp=true
RemainAfterExit=true
TimeoutStopSec=infinity

[Install]
WantedBy=multi-user.target
